# market-mindset.service
[Unit]
Description=Market Mindset Gunicorn Service
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/app
Environment="PATH=/home/ubuntu/app/venv/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/ubuntu/app/.env

# Gunicorn command
ExecStart=/home/ubuntu/app/venv/bin/gunicorn \
    --worker-class gthread \
    --workers 4 \
    --threads 4 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --graceful-timeout 60 \
    --keep-alive 5 \
    --worker-tmp-dir /dev/shm \
    --bind 127.0.0.1:8050 \
    --timeout 120 \
    --access-logfile /home/ubuntu/app/logs/access.log \
    --error-logfile /home/ubuntu/app/logs/error.log \
    --access-logformat '%(h)s %(l)s %(u)s %(t)s "%(r)s" %(s)s %(b)s "%(f)s" %(D)s' \
    --capture-output \
    application:application

# Restart policy
Restart=always
RestartSec=5
TimeoutStartSec=180
TimeoutStopSec=60

# Resource limits (optional, for stability)
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
